# s6/jenkins/run
#!/usr/bin/env bash
#
# Jenkins s6 run script
#
/usr/sbin/runuser -l  jenkins -c '/usr/bin/java -jar /usr/share/java/jenkins/jenkins.war'